# mcp-server/scripts/com.datacore.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.datacore.daily-sync</string>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/291928k/Library/CloudStorage/OneDrive-Curtin/Developer/datacore/mcp-server/scripts/daily-sync.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/291928k/Library/CloudStorage/OneDrive-Curtin/Developer/datacore/mcp-server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>DATABRICKS_HOST</key>
    <string>https://adb-7405608864534253.13.azuredatabricks.net</string>
    <key>DATABRICKS_TOKEN</key>
    <string>SET_IN_INSTALLED_COPY</string>
    <key>DATACORE_BRONZE_DIR</key>
    <string>/Users/291928k/.datacore/bronze</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/291928k/.datacore/logs/daily-sync.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/291928k/.datacore/logs/daily-sync.stderr.log</string>
</dict>
</plist>
