# mcp-server/scripts/com.datacore.gemini-session-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.datacore.gemini-session-watcher</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/291928k/david/Developer/datacore/mcp-server/scripts/gemini-session-watcher.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/291928k/david/Developer/datacore/mcp-server</string>
  <key>StandardOutPath</key>
  <string>/Users/291928k/.datacore/logs/gemini.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/291928k/.datacore/logs/gemini.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/291928k</string>
  </dict>
</dict>
</plist>
